# services.d/api/run
#!/command/with-contenv sh

cd /srv/api
exec /srv/api/containerhub-api -c /srv/api/data/config.toml
